# deployment/webhook/deploy-dev-webhook.service
[Unit]
Description=Tu-Link Webhook Deploy Listener
Documentation=https://github.com/adnanh/webhook
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/webhook -hooks /opt/tulink-backend/deployment/webhook/hooks.json -port 9000 -verbose -hotreload
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
